# rexec_server/rexec-server@.service
[Unit]
Description=rexec-server

[Service]
Type=simple
ExecStart=/usr/bin/rexec_server
User=<user>
Group=<group>
